# wan_tmb/docker/cron/crontab
* * * * * php /src/artisan schedule:run >> /dev/stdout 2>&1